# run_scrapers.cron
0 */12 * * * node scrapers/subtext.js
0 */12 * * * node scrapers/monogram.js
